# ops/openclaw/systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway for AfroTools
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/afrotools
Environment=HOME=/home/openclaw
EnvironmentFile=-/etc/openclaw/openclaw.env
ExecStart=/usr/local/bin/openclaw gateway --port 18789
Restart=always
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
